# contractor-webui-api/systemd/contractor-webui-api@.service
[Unit]
Description=contractor-webui-api (%i)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=%h/.config/contractor-webui-api/%i.env
WorkingDirectory=%h/contractor-webui-api
ExecStart=%h/contractor-webui-api/.venv/bin/uvicorn app.main:app --host ${BIND_HOST} --port ${BIND_PORT} --log-level info
Restart=on-failure
RestartSec=5
# Listen on localhost only. The CF tunnel is the public surface.
# NOTE: this unit runs as the invoking user (systemd --user). Do not enable as
# a root/system unit: it must run under the contractor's Linux user so its
# process tree + cache directory respect the per-contractor silo boundary.
NoNewPrivileges=true
PrivateTmp=false
ProtectSystem=full
ProtectHome=false

[Install]
WantedBy=default.target
